# systemd_config/caps-map-rs.service
[Unit]
Description=caps-maps-rs daemon

[Service]
Type=simple
ExecStart=/usr/bin/caps-map-rs

[Install]
WantedBy=multi-user.target


